# setup/90-can0.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="can0", RUN+="/opt/setup_can.sh"s